# meta-custom/recipes-core/my-c-app/files/my-c-app.service
[Unit]
Description=My Custom C Application
After=network.target

[Service]
ExecStart=/usr/bin/my-c-app
Restart=always
User=root

[Install]
WantedBy=multi-user.target
